# CentOS 7.x设置自定义开机启动
首先先建一个.sh的脚本文件，然后根据以下步骤进行：
1.自己定义脚本格式如下： 
必须添加的 
#!/bin/sh 
#chkconfig: 2345 80 80 
#description: auto start web server 

第一句：必须的 
声明sh 
第二句：必须的 
2345 表示linux 有 0-6种状态 在2345下启动 80：表示启动顺序 80：表示关闭顺序 
后面两个参数小于 100 
第三句：必须的 
描述 

下面可以加自己的脚本 

把定要好的脚本赋可执行权限 

2.chmod +x <自己的脚本> 

3.copy 到 /etc/init.d/ 
Sh代码 
cp <自己的脚本> /etc/init.d/  

chkconfig --add <脚本名> // 添加 
chkconfig  --list <脚本名> // 查看 

