# How to Setup a Multi-Protocol VPN Server Using SoftEther
yum install lynx -y
lynx http://www.softether-download.com/files/softether/
yum groupinstall "Development Tools"
tar xzvf sof*
cd vpnserver
make

cd ..
mv vpnserver /usr/local
cd /usr/local/vpnserver/

chmod 600 *
chmod 700 vpnserver
chmod 700 vpncmd

vi /etc/init.d/vpnserver
cat >>/etc/ipsec.conf<<EOF
#!/bin/sh
# chkconfig: 2345 99 01
# description: SoftEther VPN Server
DAEMON=/usr/local/vpnserver/vpnserver
LOCK=/var/lock/subsys/vpnserver
test -x $DAEMON || exit 0
case "$1" in
start)
$DAEMON start
touch $LOCK
;;
stop)
$DAEMON stop
rm $LOCK
;;
restart)
$DAEMON stop
sleep 3
$DAEMON start
;;
*)
echo "Usage: $0 {start|stop|restart}"
exit 1
esac
exit 0
EOF
mkdir /var/lock/subsys
chmod 755 /etc/init.d/vpnserver && /etc/init.d/vpnserver start
chkconfig --add vpnserver 
